# packaging/ferromix2.service
[Unit]
Description=FerroMix2 virtual audio mixer daemon
After=pipewire.service wireplumber.service
Wants=pipewire.service

[Service]
ExecStart=/usr/bin/ferromix2-daemon
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
